# restart_sheet.service
[Unit]
Description = Python Script
After = multi-user.target

[Service]
Type = oneshot
WorkingDirectory = /home/nvidia/inv/
ExecStart = /usr/bin/python3 sheet.py
Restart = always

[Install]
WantedBy = multi-user.target
